# src/files/docker/base-images/alpine_3_7/fluxbox
#!/sbin/openrc-run

name=fluxbox
command="/usr/bin/fluxbox"
command_args=""
command_user="abstruse"
pidfile="/var/run/fluxbox.pid"
command_background="yes"

depend() {
  need xvfb
}

start() {
  start-stop-daemon --start --quiet --pidfile $pidfile --make-pidfile --background --exec $command
}
