# src/files/docker/base-images/alpine_3_7/sshd
#!/sbin/openrc-run

name=sshd
command="/usr/sbin/sshd"
command_args="-D"
command_user="root"
pidfile="/var/run/sshd.pid"
command_background="yes"

start() {
  start-stop-daemon --start --quiet --pidfile $pidfile --make-pidfile --background --exec $command -- $command_args
}
